# supervisor/sync.ini
[program:sync]
command=/usr/bin/gunicorn --workers=2 --bind=192.168.43.106:5016 data_sync:app 
directory=/root/flask_es/
autorestart=true
autostart=true
stderr_logfile=/var/log/supervisor_flask_sync.log
user=root
startsecs=1
