# supervisor/wbsearch.ini
[program:wbsearch]
command=/usr/bin/gunicorn --workers=2 --bind=192.168.43.106:5000 app:app
directory=/root/flask_es/
autorestart=true
autostart=true
stderr_logfile=/var/log/supervisor.log
